# package/hello/etc/cron.d/test_cron_job
*/5 * * * * ads-rnr python /home/ads-rnr/main_script/hello.py