# package/hello/etc/cron.d/test_cron_job
*/1 * * * * ayush.bhatnagar python /home/ayush.bhatnagar/main_script/hello.py >> /home/ayush.bhatnagar/dummy_file