# package/hello/etc/cron.d/test_cron_job
*/1 * * * * ayush.bhatnagar bash /home/ayush.bhatnagar/main_script/run.sh
#always end this file with a new line character....
